# system/systemd/slimta-docker-check-certs.service

[Service]
Type=oneshot
ExecStart={{DIR}}/bin/check-certs {{DIR}}/.env
