# system/systemd/slimta-docker-watch-certs.path

[Unit]
BindsTo=slimta-docker.service

[Path]
PathChanged=/etc/ssl/private/local/
