# system/systemd/slimta-docker-check-certs.timer

[Unit]
BindsTo=slimta-docker.service

[Timer]
OnCalendar=daily EST
RandomizedDelaySec=3600
Persistent=true
