# systemd/slimta-docker-check-certs.service

[Unit]
ConditionPathExists=/etc/letsencrypt.env

[Service]
Type=oneshot
ExecStart=/opt/slimta-docker/bin/check-certs
