# systemd/slimta-docker-watch-certs.service

[Service]
Type=oneshot
ExecStart=docker-compose restart pymap slimta-edge slimta-relay
WorkingDirectory=/opt/slimta-docker/
